# supabase/functions/send-reminders/crontab
*/5 * * * * curl -i -X POST https://<project-ref>.supabase.co/functions/v1/send-reminders